# config/systemd-templates/user/pycardsyncer.service.template
[Unit]
Description=PyCardDav Syncer

[Service]
Type=oneshot
ExecStart=/home/kenyonj/.pyenv/shims/pycardsyncer
