# config/systemd-templates/user/spotifysync.timer.template
[Unit]
Description=Spotify Playlist Syncer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
