# config/systemd-templates/user/spotifysync.service.template
[Unit]
Description=Spotify Playlist Syncer

[Service]
Type=oneshot
ExecStart=/home/kenyonj/.scripts/spotifysyncer
