# tag-music/config/systemd/user/spotifysync.service
[Unit]
Description=Spotify Playlist Syncer

[Service]
Type=oneshot
ExecStart=$HOME/.scripts/spotifysyncer

[Install]
WantedBy=default.target
