# tag-mail-receiver/config/systemd/user/offlineimap.timer
[Unit]
Description=offlineimap timer

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
