# tag-mail-sender/config/systemd/user/pycardsyncer.service
[Unit]
Description=PyCardDav Syncer

[Service]
Type=oneshot
ExecStart=/usr/bin/pycardsyncer
